# deploy/systemd/tradebot-news.timer
[Unit]
Description=Run the tradebot XSP/MCL news inference approximately every four hours

[Timer]
OnBootSec=15min
OnUnitInactiveSec=4h
AccuracySec=15min
Unit=tradebot-news.service

[Install]
WantedBy=timers.target
